# login/brain-launcher.service
[Unit]
Description=BrAIn-nw1 Launcher App
After=network.target cloudflared-brain.service
Wants=cloudflared-brain.service

[Service]
Type=simple
User=he
WorkingDirectory=/home/he/ai/dev/langgraph/KB_BS_local-rag-he/login
Environment="LAUNCHER_PASSWORD=BrAIn-#1"
Environment="LAUNCHER_URL=https://brain-nw1-launcher.trycloudflare.com"
Environment="MAIN_APP_URL=https://brain-nw1-app.trycloudflare.com"
Environment="PATH=/home/he/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/he/.local/bin/uv run --with streamlit --with psutil streamlit run launcher_app.py --server.port 8502 --server.headless true
Restart=always
RestartSec=10s
StandardOutput=append:/tmp/launcher.log
StandardError=append:/tmp/launcher.log

[Install]
WantedBy=multi-user.target
