# login/systemd-template.service
[Unit]
Description=RAG Researcher Launcher
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/ai/dev/langgraph/KB_BS_local-rag-he/login
Environment="LAUNCHER_PASSWORD=YOUR_SECURE_PASSWORD"
ExecStart=/usr/bin/streamlit run launcher_app.py --server.port 8502 --server.headless true
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
